# containers/apache/2/php/5.6/rootfs/etc/cron.d/default
# События Bitrix
* * * * * php -f /var/www/html/local/php_interface/cron_events.php >> /var/log/app/cron.log
